# crontab
# Run the webscraper at 12:00 and 19:00 every day
10 12,19 * * * cd ~/Dev/MakelaarWebscraper/Makelaar-webscraper && python huurhuis_webscraper.py >> logs/cron.log 2>&1
