# crontab
# Run the webscraper at 12:00 and 18:00 every day
0 12,19 * * * cd /app && python /app/huurhuis_webscraper.py >> /app/logs/cron.log 2>&1
